# etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/main
::restart:-/bin/sh
::ctrlaltdel:/bin/reboot
::shutdown:/bin/umount -a -r
